# systemd/backup_postgresql.timer
[Unit]
Description=Run backups on PostgreSQL databases

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
